# systemd/github-auto-push.timer
[Unit]
Description=Weekly GitHub auto-push for claude_discord

[Timer]
OnCalendar=Mon *-*-* 09:17:00
Persistent=true

[Install]
WantedBy=timers.target
