# systemd/multi-repo-sync.timer
[Unit]
Description=Weekly multi-repo GitHub auto-sync

[Timer]
OnCalendar=Mon *-*-* 09:17:00
Persistent=true

[Install]
WantedBy=timers.target
